# misc/cron.d/healthcheck
#
# cron-jobs for healthcheck tripwire
#

MAILTO=root

*/10 * * * *	root	python /opt/dvmps/bin/healthcheck.py

